# compose/django/supervisord.conf
[program:dc_traffic_tickets]
command=docker-compose up
redirect_stderr=true
autostart=true
autorestart=true
priority=10
